# config/supervisor/image-master-push-worker.conf
; Image Master marketplace push worker
; Copy to /etc/supervisor/conf.d/image-master-push-worker.conf on server
; Run: sudo supervisorctl reread && sudo supervisorctl update && sudo supervisorctl start image-master-push-worker:*

[program:image-master-push-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/inventory_5c_usr/data/www/inventory.5coremanagement.com/artisan queue:work database --queue=image-master-push --sleep=3 --tries=1 --timeout=14400 --max-time=14400
directory=/var/www/inventory_5c_usr/data/www/inventory.5coremanagement.com
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=inventory_5c_usr
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/inventory_5c_usr/data/www/inventory.5coremanagement.com/storage/logs/image-master-push-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=14400
stopsignal=SIGTERM
